# deploy/superset.service
# Systemd service file for Superset
#
# This file allows you to run the Superset application as a systemd service.
#
# Before using, make sure the paths in this file are correct for your system.
#
# Installation:
# 1. sudo cp deploy/superset.service /etc/systemd/system/
# 2. sudo systemctl daemon-reload
# 3. sudo systemctl enable superset
# 4. sudo systemctl start superset
#

[Unit]
Description=Superset service for FormaSup BI Platform
After=docker.service
Requires=docker.service

[Service]
User=your_user
Group=your_group
WorkingDirectory=/path/to/your/project/formasup-bi-platform
Restart=always

# Compose up
ExecStart=/usr/bin/docker-compose -f docker-compose.yml -f docker-compose.prod.yml up
# Compose down
ExecStop=/usr/bin/docker-compose -f docker-compose.yml -f docker-compose.prod.yml down

[Install]
WantedBy=multi-user.target